# deployment/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/local/bin/gunicorn application:application -w 4 --config /src/gunicorn.conf --log-config /src/logging.conf --bind 0.0.0.0:5000
directory=/src
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler